# .config/profile.d/asdf.sh
export ASDF_GEM_DEFAULT_PACKAGES_FILE="$HOME/.config/asdf/default-gems"
export ASDF_RUBY_BUILD_VERSION=master
