# .config/profile.d/brew.sh
type brew > /dev/null 2>&1 || return

export HOMEBREW_NO_ANALYTICS=1
export HOMEBREW_NO_AUTO_UPDATE=1
export HOMEBREW_NO_INSTALL_CLEANUP=1
