# .config/profile.d/coreutils.sh
export INPUTRC=~/.config/coreutils/inputrc
