# .config/profile.d/readline.sh
export INPUTRC=~/.config/readline/inputrc
