# .config/profile.d/vim.sh
export VIMINIT='source ~/.config/vim/vimrc'
